# scripts/cron/reference_sync.cron
# Reference-data sync (Option C). Install with: crontab scripts/cron/reference_sync.cron
# Adjust cadence to how often the reference tables change (daily here, 02:15 UTC).
15 2 * * * cd /home/ec2-user/CTV_dbt_iceberg_poc && docker compose exec -T ingestion python -m ingestion.reference_sync >> /var/log/ref_sync.log 2>&1
